# Daemon/backup_daemon.service
[Unit]
Description=Backup Daemon

[Service]
ExecStart=/home/user/daemon/main.cpp
Restart=always

[Install]
WantedBy=default.target

